# tools/native/openwrt/files/etc/uci-defaults/99-trustforge
#!/bin/sh
# SPDX-License-Identifier: Apache-2.0 OR MIT
#
# uci-defaults runs once on first boot after package install. It is the
# canonical place to seed defaults that depend on the live system
# (hostname, MAC, etc) rather than fixed values shipped in /etc/config.
#
# This script:
#   1. ensures a `daemon main` UCI section exists,
#   2. seeds a stable instance id derived from the device hostname,
#   3. ensures the runtime directory and ownership are correct,
#   4. enables and starts the service.
#
# It removes itself on success so it never runs twice.

set -e

CONF=trustforge

. /lib/functions.sh

uci -q get $CONF.main >/dev/null || {
	uci batch <<-EOF
		set $CONF.main=daemon
		set $CONF.main.enabled='1'
		set $CONF.main.listen='127.0.0.1:8642'
		set $CONF.main.socket='/var/run/trustforge/decide.sock'
		set $CONF.main.profile='tf-home-compatible'
	EOF
}

# Stable per-device instance id under tf:instance:device:<hostname>/openwrt
HOSTNAME=$(uci -q get system.@system[0].hostname || echo "openwrt")
uci set $CONF.main.instance_id="tf:instance:device:${HOSTNAME}/openwrt"

uci commit $CONF

mkdir -p /var/run/trustforge /etc/trustforge
chown trustforge:trustforge /var/run/trustforge 2>/dev/null || true

# Materialise a minimal /etc/trustforge/config.yaml if none exists so
# tf-daemon does not refuse to start on first boot.
if [ ! -f /etc/trustforge/config.yaml ]; then
	cat > /etc/trustforge/config.yaml <<-EOF
		# Auto-generated by 99-trustforge on first boot.
		# Edit /etc/config/trustforge (UCI) for persistent changes.
		profile: tf-home-compatible
		log_level: info
	EOF
	chown -R trustforge:trustforge /etc/trustforge
fi

/etc/init.d/trustforge enable
/etc/init.d/trustforge start || true

exit 0
